# autostart/crontab
@reboot sleep 30 && export DISPLAY=:0 && sh ~/keypress.sh