# lattepanda M2 tester/crontab
@reboot sleep 13 && export DISPLAY=:0 && sh ~/M2Tester.sh

